# etc/init/distributor-agent.conf
description "Start up the Octavia Distributor Agent"

start on startup

respawn
respawn limit 2 2

exec distributor-agent --config-file /etc/octavia/distributor-agent.conf
